# .config/autostart/init.desktop
[Desktop Entry]
Comment=My personal scripts
Exec=/home/user/.config/init.sh
Name=Init.sh
Type=Application
